# rootfs/etc/s6-overlay/s6-rc.d/cli-setup/up
/bin/bash -c "mkdir -p /home/ubuntu/.local/bin /home/ubuntu/.opencode/bin && if [ -f /usr/local/bin/claude ]; then mv /usr/local/bin/claude /home/ubuntu/.local/bin/claude && chown ubuntu:ubuntu /home/ubuntu/.local/bin/claude; fi && if [ -f /usr/local/bin/opencode ]; then mv /usr/local/bin/opencode /home/ubuntu/.opencode/bin/opencode && chown ubuntu:ubuntu /home/ubuntu/.opencode/bin/opencode; fi"
