# rootfs/etc/s6-overlay/s6-rc.d/ttyd-codex/run
#!/command/with-contenv bash
cd /home/ubuntu
export HOME=/home/ubuntu
exec s6-setuidgid ubuntu \
  ttyd --interface ${BIND_ADDR:-0.0.0.0} --port 7685 --writable \
    --credential "${TTYD_CREDENTIAL}" \
    codex
